# bin/massema.net.service
[Unit]
Description=massema.net
After=network.target

[Service]
User=luser
ExecStart=/usr/bin/massema.net foreground
ExecRestart=/usr/bin/massema.net reboot
ExecStop=/usr/bin/massema.net stop
Restart=on-failure

[Install]
WantedBy=multi-user.target