# bin/massema.net
#!/bin/bash
### BEGIN INIT INFO
# Provides:           massema.net
# Required-Start:     $remote_fs $network
# Required-Stop:      $remote_fs $network
# Default-Start:      2 3 4 5
# Default-Stop:       0 1 6
# Short-Description:  Start massema.net
### END INIT INFO
#
# the below requires a somewhat debian-like system.
#
# make a link to this file from etc/init.d;
# ln -s <this file> /etc/init.d/
#
# create the runlevel-links like this;
# sudo update-rc.d massema.net defaults 50
#
# now this should work;
# invoke-rc.d massema.net start
# invoke-rc.d massema.net stop

# $1 is a list of dirs. return the first dir that is writable.
tryt() {
    while [ 0 -lt $# ]; do
        if 2>/dev/null touch $1/$$ && rm $1/$$; then
            echo $1
            return
        else
            shift
        fi
    done
    exit -1
}

export ITEM=massema.net
export ITEM_STARTMOD=mm
export ITEM_USER=masse
export ITEM_LOGDIR=`tryt /var/log /tmp`/$ITEM
export ITEM_ROOT=`eval echo ~$ITEM_USER`/git/$ITEM
export ITEM_LIB="$ITEM_ROOT/_build/default/lib"
export ITEM_ERL=`which erl`
export ITEM_GROUP=`id -gn $ITEM_USER`
export ITEM_BOOTLOG=$ITEM_LOGDIR/boot.log
export ITEM_ERLLOG=$ITEM_LOGDIR/erlang.log
export ITEM_SNAME=`echo $ITEM | tr "\-." "_"`

[ -d $ITEM_LOGDIR ] || mkdir $ITEM_LOGDIR
chown $ITEM_USER:$ITEM_GROUP $ITEM_LOGDIR
[ -f $ITEM_BOOTLOG ] || touch $ITEM_BOOTLOG
chown $ITEM_USER:$ITEM_GROUP $ITEM_BOOTLOG

# osx does not have pkill. #fail
beamgrep() {
    ps -u "$1" -o pid,command | grep beam | grep "$2" | cut -c-6
}

item_start() {
    P=`beamgrep $ITEM_USER "run $ITEM_STARTMOD"`
    [ -n "$P" ] && return
    for l in $ITEM_LIB/*; do
        [ -d $l/ebin ] && PAS="$PAS -pa $l/ebin"
    done
    [ -z "$ITEM_ERL" ] && echo "No erl in path" && exit -1
    sudo -u $ITEM_USER \
        $ITEM_ERL  \
        -sname $ITEM_SNAME \
        -setcookie $ITEM \
        -boot start_sasl \
        -kernel error_logger "{file,\"$ITEM_ERLLOG\"}" \
        $PAS \
        -run $ITEM_STARTMOD \
        -detached
}

item_stop() {
    P=`beamgrep $ITEM_USER "run $ITEM_STARTMOD"`
    [ -n "$P" ] && sudo -u $ITEM_USER kill -KILL $P
}

item_shell() {
    item_remsh
}

item_remsh() {
    HOST=`echo $HOSTNAME | cut -f1 -d"."`
    sudo -u $ITEM_USER \
        $ITEM_ERL \
        -sname $$ -remsh ${ITEM_SNAME}@${HOST} -setcookie $ITEM $*
}

out() {
    flag=""
    [ "$1" == "-n" ] && flag="-n" && shift
    echo $flag "$1" >> $ITEM_BOOTLOG
}

otp_release() {
    $ITEM_ERL -noshell \
        -eval 'io:format("~s~n", [erlang:system_info(otp_release)])' \
        -s erlang halt
}

out "`date | tr ' ' '-'` "
case "$1" in
    status)
        P=`beamgrep $ITEM_USER "run $ITEM_STARTMOD"`
        if [ -n "$P" ] ; then
            out "status: running"
            echo "pid       : $P"
            echo "logdir    : $ITEM_LOGDIR"
            echo "rootdir   : $ITEM_ROOT"
            echo "erl       : $ITEM_ERL ("`otp_release`")"
            echo "running as: $ITEM_USER:$ITEM_GROUP"
        else
            out "status: halted"
            echo "not started"
        fi
        ;;
    start)
        out "logging to "$ITEM_LOGDIR
        out -n "starting $0..."
        item_start && out -n " started..."
        out " done."
        ;;
    stop)
        out -n "stopping..."
        item_stop && out -n " stopped..."
        out " done."
        ;;
    restart | force-reload)
        out "restarting."
        $0 stop
        sleep 1
        $0 start
        ;;
    attach | shell)
        out "attaching $USER - $$."
        item_shell || out -n "Failed... "
        out "detaching $USER - $$."
        ;;
    reload-geodata)
        $ITEM_LIB/$ITEM_STARTMOD/priv/bin/reload_geodata $ITEM $ITEM_SNAME
        ;;
    *)
        echo "argument '"$1"' unknown."
        exit 1
esac
